# package/ddns-gargoyle/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 84
/ericpaulbishop/gargoyle.git/!svn/ver/3378/package/ddns-gargoyle/files/etc/hotplug.d
END
